# systemd/agentverif-api.service
[Unit]
Description=agentverif API
Documentation=https://agentverif.com/docs
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=/root/agentverif/api
ExecStart=/root/agentverif/.venv/bin/uvicorn server:app --host 0.0.0.0 --port 8090 --workers 2
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agentverif-api

# Environment
Environment=AGENTVERIF_DB=/root/agentverif/api/agentverif.db
Environment=AGENTVERIF_API_KEY=5be8d4bbd1727c127cbccbc020664f5d244e0519b8fb61daae4423a5dcb36f21
Environment=AGENTVERIF_SCAN_URL=https://api.agentverif.com/scan

# Security hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
